# pi/farmon.service
[Unit]
Description=FarMon IoT Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Bind to all interfaces on port 8090 so LAN devices and browsers can reach the UI/API.
# Change 0.0.0.0 to 127.0.0.1 if you only want local access.
ExecStart=/opt/farmon/farmon serve --http=0.0.0.0:8090
WorkingDirectory=/opt/farmon
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/farmon/pb_data

[Install]
WantedBy=multi-user.target
